# roles/domain/files/etc/cron.d/domain
# shell variable for cron
SHELL=/bin/bash
# PATH variable for cron
PATH=/usr/local/bin:/usr/local/sbin:/sbin:/usr/sbin:/bin:/usr/bin
#M   S     T M W  User  Befehl
#-----------------------------------------------------------------
39    12     * * *  root  python -m cryptdomainmgr --next /etc/cryptdomainmgr/domains.conf /etc/cryptdomainmgr/inwxcred.conf

